# playbooks/roles/rsyslog_config/templates/ovhcloud-sap-audit.service.j2
[Service]
Environment="OVHCLOUD_SAP_SID={{ sap_rsyslog_sid }}"
Environment="OVHCLOUD_SAP_INSTANCE_NUMBER={{ instno_service }}"
Environment="OVHCLOUD_SAP_SOFTWARE_STACK={{ software_stack }}"